# roles/factorio/templates/factorio-s3-sync-systemd.service.j2
[Unit]
Description=Factorio S3 Sync Service
Wants=factorio.service
After=factorio.service

[Service]
Type=simple
User=factorio
ExecStart={{ factorio_home_dir }}/s3_sync_watch.sh
